# recipes-graphics/drivers/visionfive2-pvr-graphics/rc.pvr.service
[Unit]
Description=PowerVR consumer services

[Service]
ExecStart=/usr/bin/rc.pvr start
ExecStop=/usr/bin/rc.pvr stop
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
