# recipes-bsp/k3-rootfs-overlay/files/weston/etc/systemd/system/weston.service.d/autodetect.conf
[Service]
ExecStart=
ExecStart=/usr/local/bin/weston-autostart.sh
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60
